# data/com.chrisdaggas.control-center-snapshot.timer
[Unit]
Description=Weekly auto-snapshot for Control Center
Documentation=https://chrisdaggas.com

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
